# media-automount@.service
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/media-automount %I
ExecStop=/usr/bin/media-automount %I
